# systemd/mcp-splitwise.service
[Unit]
Description=MCP Splitwise Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/assistant/mcp-servers/splitwise-mcp-server
EnvironmentFile=/home/ubuntu/assistant/.env.production
ExecStart=/usr/bin/node dist/index.js --http --port 4000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-splitwise

[Install]
WantedBy=multi-user.target
