# systemd/mcp-google-places.service
[Unit]
Description=MCP Google Places Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/assistant/mcp-servers/googleplaces-mcp-server
EnvironmentFile=/home/ubuntu/assistant/.env.production
ExecStart=/usr/bin/node dist/index.js --http --port 1111
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-google-places

[Install]
WantedBy=multi-user.target
